# hadoop-master1/files/etc/supervisord.d/ranger-admin.conf
[program:ranger-admin]
directory=/etc/ranger
command=bash startup.sh
startsecs=30
stopwaitsecs=10
user=ranger
redirect_stderr=true
stdout_logfile=/var/log/ranger-admin.log
autostart=true
