# base-ubuntu-2204/files/etc/supervisord.d/sshd.conf
[program:sshd]
directory=/usr/sbin
command=/usr/sbin/sshd -D
startsecs=30
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/sshd.log
autostart=true
autorestart=true