# hadoop-worker2/files/etc/supervisord.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
