# templates/hadoop-master/files/etc/supervisor/conf.d/hive-server2.conf
[program:hive-server2]
directory=/opt/hive
command=/opt/start-hive-server2.sh
startsecs=30
stopwaitsecs=10
user=hive
redirect_stderr=true
stdout_logfile=/var/log/hive-server2.log
autostart=true
