# templates/hadoop-master/files/etc/supervisord.d/hive-server2.conf
[program:hive-server2]
directory=/opt/hive
command=hive --service hiveserver2
startsecs=30
stopwaitsecs=10
user=hive
redirect_stderr=true
stdout_logfile=/var/log/hive-server2.log
autostart=true
