# templates/hadoop-master/files/etc/supervisord.d/kyuubi-server.conf
[program:kyuubi-server]
directory=/opt/kyuubi
command=kyuubi run
startsecs=30
stopwaitsecs=10
user=kyuubi
redirect_stderr=true
stdout_logfile=/var/log/kyuubi-server.log
autostart=true
